# systemd/raspi-tailnet-monitor.service
[Unit]
Description=Raspberry Pi Tailnet Monitor
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/raspi-tailnet-monitor
ExecStart=/usr/bin/python3 /opt/raspi-tailnet-monitor/app.py --config /etc/raspi-tailnet-monitor/config.json
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
# Lock down the home directory except for the SSH key directory (read for
# authentication) and the SSH multiplexing control socket directory (read/write
# for connection reuse across poll cycles).
ProtectHome=read-only
BindReadOnlyPaths=/home/pi/.ssh
BindReadWritePaths=/home/pi/.cache/raspi-tailnet-monitor
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
