# systemd/raspi-tailnet-monitor.service
[Unit]
Description=Raspberry Pi Tailnet Monitor
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/raspi-tailnet-monitor
ExecStart=/usr/bin/python3 /opt/raspi-tailnet-monitor/app.py --config /etc/raspi-tailnet-monitor/config.json
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
# Lock down the home directory except for the SSH key directory, which the
# monitor needs read access to in order to authenticate to remote devices.
ProtectHome=read-only
BindReadOnlyPaths=/home/pi/.ssh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
